# healthcoach/rootfs/etc/services.d/mcp/finish
#!/command/with-contenv bash
# Delay before s6 restarts the service (prevents tight restart loops).
sleep 5
