# healthcoach/rootfs/etc/services.d/sync/run
#!/command/with-contenv bash
# Async task runner (APScheduler) — periodic data sync + PMC calculation.
cd /app
exec python3 /app/scripts/task_runner.py
